# templates/centos/ceph-rest-api.service.erb
[Unit]
Description=Ceph RESTful API for admin purposes
After=network-online.target local-fs.target
Wants=network-online.target local-fs.target

[Service]
Type=Notify
Environment=CLUSTER=<%= node['ceph']['cluster'] %>
ExecStart=/usr/bin/python /bin/ceph-rest-api --cluster ${CLUSTER} -c /etc/ceph/${CLUSTER}.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
